# sound.service
[Unit]
Description=Sound

[Service]
ExecStart=/usr/bin/Sound.sh

[Install]
WantedBy=multi-user.target
